# rc.local
## remove any old lock files
rm -rf /zeus/controller/*.loc

exit 0
